# ubuntu-server/ddos-dashboard.service
[Unit]
Description=Enterprise DDoS Web Management Dashboard
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/ddos-enterprise-lab/ubuntu-server/web-dashboard
ExecStart=/usr/bin/python3 /root/ddos-enterprise-lab/ubuntu-server/web-dashboard/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
